# systemd/hs100@.timer
[Unit]
Description=Periodic checking of %I status

[Timer]
OnUnitActiveSec=5s
AccuracySec=1
OnBootSec=10s

[Install]
WantedBy=hs100@%i.service